# deploy/backup/crontab
30 2 * * * docker compose -f /work/docker-compose.yml --profile backup run --rm backup >> /var/log/backup.log 2>&1
